# backend/control.service
[Unit]
Description=Senior Design Node Control Service
After=network.target

[Service]
Type=simple
User=jbardwick
WorkingDirectory=/home/jbardwick/Desktop
Environment="PATH=/home/jbardwick/Desktop/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Clean up any left-over processes
ExecStartPre=/usr/bin/pkill -f "python3 receive.py"
ExecStartPre=/usr/bin/pkill -f "npm"
# Run a shell command that performs:
# 1. Activation of the virtual environment.
# 2. Starting the frontend and backend (using npm start) in their respective directories.
# 3. Starting the receive.py script.
ExecStart=/bin/bash -c ' \
  cd /home/jbardwick/Desktop && \
  source venv/bin/activate && \
  cd senior-design-UI/backend && npm start & \
  cd ../frontend && npm start & \
  cd ../messaging && python3 receive.py'
Restart=always
RestartSec=10
KillMode=control-group
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
